# build/config_files/supervisord.conf
[supervisord]
nodaemon=true

[program:startup.sh]
command = bash -x /var/www/html/startup.sh
startsecs = 0
autorestart = false
startretries = 1

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true
